# packages/orbit-cli/orbit_cli-0.1.26-py3-none-any.whl/orbit_cli/templates/application/apt/DEBIAN/postinst
#!/bin/bash

adduser --system --disabled-login --quiet --shell /bin/bash --gecos "{{ project }}" {{ project }} --home /home/{{ project }}
systemctl enable {{ project }}
service {{ project }} stop 2> /dev/null
sleep 1
service {{ project }} start
