# packages/matrix-reportbot/matrix_reportbot-0.0.2.tar.gz/matrix_reportbot-0.0.2/reportbot.service
[Unit]
Description=ReportBot
Requires=network.target

[Service]
Type=simple
User=reportbot
Group=reportbot
WorkingDirectory=/opt/reportbot
ExecStart=/opt/reportbot/.venv/bin/python3 -um matrix_reportbot
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target