# packages/orbit-cli/orbit_cli-0.1.26.tar.gz/orbit_cli-0.1.26/orbit_cli/templates/application/apt/etc/systemd/system/project.service
[Unit]
Description={{ project }} Daemon service
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/home/{{ project }}
ExecStart=/opt/{{ project }}/{{ project }} --run
User={{ project }}
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
